# db-replica/config/db-backup.service
[Unit]
Description="MySQL DB backup"

[Service]
ExecStart=/usr/local/bin/db-backup.sh
